# src/variants/v8/filesystem/root/etc/rc.local
xinit ./etc/directly_to_browser.sh
exit 0
